# php/var/spool/cron/crontabs/root
# m h  dom mon dow   command
*/5 *   * * *   php /var/www/pussybnb.com/wp-cron.php
